# monitor-and-notify.service
[Unit]
Description=My Bluetooth Notification Service
After=multi-user.target
Wants=minute-timer.timer

[Timer]
OnCalendar=*:0/1

[Service]
Type=idle
User=pi
WorkingDirectory=/home/pi/Documents/workspace/TemperatureMonitoringSystem/Documents/workspace/TemperatureMonitoringSystem/Documents/workspace/TemperatureMonitoringSystem/Documents/workspace/TemperatureMonitoringSystem/Documents/workspace/TemperatureMonitoringSystem/Documents/workspace/TemperatureMonitoringSystem/Documents/workspace/TemperatureMonitoringSystem
ExecStart=/usr/bin/python3 /home/pi/Documents/workspace/TemperatureMonitoringSystem/Documents/workspace/TemperatureMonitoringSystem/Documents/workspace/TemperatureMonitoringSystem/Documents/workspace/TemperatureMonitoringSystem/Documents/workspace/TemperatureMonitoringSystem/Documents/workspace/TemperatureMonitoringSystem/Documents/workspace/TemperatureMonitoringSystem/monitorAndNotify.py > /home/pi/Documents/workspace/TemperatureMonitoringSystem/Documents/workspace/TemperatureMonitoringSystem/Documents/workspace/TemperatureMonitoringSystem/Documents/workspace/TemperatureMonitoringSystem/Documents/workspace/TemperatureMonitoringSystem/Documents/workspace/TemperatureMonitoringSystem/Documents/workspace/TemperatureMonitoringSystem/monitorAndNotify.log 2>&1

[Install]
WantedBy=multi-user.target
WantedBy=minute-timer.target